# crontab
3 * * * * root /usr/local/apache2/generate_vanity_map.sh
# An empty line is required at the end of this file for a valid cron file.
